# proccess.conf
[program:proccess]
command=/api/index.py
directory=/api
autostart=true
autorestart=true
stderr_logfile=/var/log/myprocess.err.log
stdout_logfile=/var/log/myprocess.out.log
