# proccess.conf
[supervisord]
logfile=/var/log/supervisord.log ; main log file; default $CWD/supervisord.log
logfile_maxbytes=50MB        ; max main logfile bytes b4 rotation; default 50MB
logfile_backups=10           ; # of main logfile backups; 0 means none, default 10
loglevel=info                ; log level; default info; others: debug,warn,trace
pidfile=/var/run/supervisord.pid ; supervisord pidfile; default supervisord.pid
nodaemon=false               ; start in foreground if true; default false
minfds=1024                  ; min. avail startup file descriptors; default 1024
minprocs=200                 ; min. avail process descriptors; default 200
user=root                    ; default is current user, required if root
childlogdir=/var/log/supervisor ; 'AUTO' child log dir, default $TEMP
environment=ENV_VAR1="value1",ENV_VAR2="value2" ; environment variables to pass to programs started by supervisord

[program:proccess]
command=/api/script.py
directory=/api
autostart=true
autorestart=true
stderr_logfile=/var/log/myprocess.err.log
stdout_logfile=/var/log/myprocess.out.log
